# lms-restart.cron
#!/bin/bash

exec /usr/bin/lms-restart >>/var/log/lms-restart.log 2>&1
